# tools/codercat-gallery.service
[Unit]
Description=Codercat Gallery
After=network.target

[Service]
Environment=S3_KEY=<key>
Environment=S3_SECRET=<secret>
Environment=S3_BUCKET=<bucket>
WorkingDirectory=/home/kiko/gallery/
Type=simple
Restart=always
RestartSec=1
User=kiko
ExecStart=/home/kiko/gallery/gallery

[Install]
WantedBy=multi-user.target